# pkg/linux/install.sh
#!/bin/bash

sudo cp xc /usr/bin
sudo cp xc.service /lib/systemd/system
sudo systemctl enable xc
sudo systemctl start xc

sudo apt update
sudo apt install -y nginx
sudo cp xc.api.pointer.pw /etc/nginx/sites-available
sudo ln -s /etc/nginx/sites-available/api.pointer.pw /etc/nginx/sites-enabled/api.pointer.pw
sudo nginx -s reload