# pkg/linux/xc.service
[Unit]
Description=XC

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/bin/xc -address=localhost:8090 -xcRoot="https://sandbox.open-xchange.com" -xcUser="user.three@sandbox-1.open-xchange.com" -xcPassword="secret"

[Install]
WantedBy=multi-user.target